# deploy/systemd/utw-match-recent.timer
[Unit]
Description=Match recently crawled tenders every hour

[Timer]
OnCalendar=hourly
Persistent=true
Unit=utw-match-recent.service

[Install]
WantedBy=timers.target
